# webssh.service
[Unit]
Description=WebSSH terminal interface
After=network.target

[Service]
User=www-data
Group=www-data
ExecStart=wssh --fbidhttp=False --port=2083

[Install]
WantedBy=multi-user.target
